# daily_log_fix_prompt_bot/systemd/daily-log-fix-prompt-bot.timer
[Unit]
Description=Daily CryptoMaster Log Analysis and Fix Prompt Bot Timer
Requires=daily-log-fix-prompt-bot.service

[Timer]
# Run daily at 08:00 UTC (shifted +2h from original 06:00 UTC schedule)
OnCalendar=*-*-* 08:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
