# systemd/cryptomaster-dashboard.service
[Unit]
Description=CryptoMaster Dashboard Web Server
After=network-online.target cryptomaster.service
Wants=network-online.target
# NOTE: intentionally NOT PartOf=cryptomaster.service. The dashboard is an
# INDEPENDENT read-only service; tying it to the bot made every bot restart
# also restart the dashboard, so a burst of bot restarts tripped the start-rate
# limit and systemd gave up — leaving the Android API down until a manual
# restart. Decoupled: bot restarts no longer touch the dashboard.
# Never give up restarting the API (belt-and-suspenders with the decoupling).
StartLimitIntervalSec=0

[Service]
Type=simple
# SHIP-DARK OPEN posture (default). This unit must start on a STOCK server with
# no extra provisioning — that is what keeps the Android API reachable after an
# autodeploy. The audit PR5 hardening (dedicated non-root user + bearer-token
# credential + read-only mount) is OPT-IN and is applied only by
# hetzner-one-time-install once the user/group, credential file, and read-grants
# exist. Installing the hardened form on an unprovisioned server makes systemd
# refuse to start (User= missing, LoadCredential= source absent) and takes the
# API down — this happened on 2026-07-19; do NOT re-add User=/LoadCredential=
# here without also provisioning them on the server.
#
# Runs as root (same as cryptomaster.service) so it can read the root-owned
# state files (cache.sqlite, learning-state JSON, open-positions JSON).
WorkingDirectory=/opt/cryptomaster
Environment="PYTHONUNBUFFERED=1"
# Bind 0.0.0.0 so the external Android app can reach the API. Auth is opt-in via
# DASHBOARD_SECURITY_ENABLED (default off = open); when it is enabled the app
# also flips its default bind to localhost. Public exposure should ultimately go
# through a VPN or authenticated HTTPS reverse proxy (audit 9.2).
Environment="DASHBOARD_BIND_HOST=0.0.0.0"
Environment="DASHBOARD_PORT=5001"

# PERMANENT FIX (CLAUDE.md 2026-06-11 + commit 78a564c 2026-07-03):
# - Direct python (no Gunicorn/WSGI wrappers).
# - Venv lives OUTSIDE the repo (/opt/dashboard_venv, built from
#   requirements-dashboard.txt) because deploy runs a hard reset.
# - start_flask_dashboard.py (tracked) fixes sys.path/cwd so
#   src.services.* imports work; do not exec dashboard_web.py directly.
ExecStart=/opt/dashboard_venv/bin/python3 -u /opt/cryptomaster/start_flask_dashboard.py
Restart=always
RestartSec=10
# (start-rate limiting is disabled via StartLimitIntervalSec=0 in [Unit] above,
# so the API always keeps retrying and never stays down after a crash burst.)

# Light, provisioning-free hardening only. The dashboard opens sqlite read-only
# (mode=ro) and reads JSON — it never writes under /opt. These directives do not
# require a dedicated user or credential and do not block root reads, so the unit
# still starts on a stock server. Heavier sandboxing (non-root User=,
# ProtectSystem=strict, ReadOnlyPaths, dropped capabilities, LoadCredential) is
# opt-in via hetzner-one-time-install after the user/credential are provisioned.
NoNewPrivileges=true
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dashboard-web

[Install]
WantedBy=multi-user.target
