# systemd/cryptomaster-dashboard.service
[Unit]
Description=CryptoMaster Dashboard Web Server
After=network-online.target cryptomaster.service
Wants=network-online.target
# NOTE: intentionally NOT PartOf=cryptomaster.service. The dashboard is an
# INDEPENDENT read-only service; tying it to the bot made every bot restart
# also restart the dashboard, so a burst of bot restarts tripped the start-rate
# limit and systemd gave up — leaving the Android API down until a manual
# restart. Decoupled: bot restarts no longer touch the dashboard.
# Never give up restarting the API (belt-and-suspenders with the decoupling).
StartLimitIntervalSec=0

[Service]
Type=simple
# Audit PR5 (P1.6): run as a dedicated non-root user with ONLY read access to
# the three data files the dashboard needs (cache.sqlite, learning-state JSON,
# open-positions JSON). Provision before deploy:
#   groupadd -r cryptomaster-dashboard
#   useradd -r -g cryptomaster-dashboard -s /usr/sbin/nologin cryptomaster-dashboard
#   # grant read on /opt/cryptomaster/{local_learning_storage,server_local_backups,data}
User=cryptomaster-dashboard
Group=cryptomaster-dashboard
WorkingDirectory=/opt/cryptomaster
Environment="PYTHONUNBUFFERED=1"
# Safe default bind: localhost only. Public exposure must go through a VPN or an
# authenticated HTTPS reverse proxy (audit 9.2) — never 0.0.0.0.
Environment="DASHBOARD_BIND_HOST=127.0.0.1"
Environment="DASHBOARD_PORT=5001"
# Bearer-token auth: systemd hands the token to the process as a credential;
# Python reads it from $CREDENTIALS_DIRECTORY/dashboard_api_token. The file is
# root-owned (0400) and never enters the repo, env, logs, or URLs.
LoadCredential=dashboard_api_token:/etc/cryptomaster/dashboard_api_token

# PERMANENT FIX (CLAUDE.md 2026-06-11 + commit 78a564c 2026-07-03):
# - Direct python (no Gunicorn/WSGI wrappers).
# - Venv lives OUTSIDE the repo (/opt/dashboard_venv, built from
#   requirements-dashboard.txt) because deploy runs a hard reset.
# - start_flask_dashboard.py (tracked) fixes sys.path/cwd so
#   src.services.* imports work; do not exec dashboard_web.py directly.
ExecStart=/opt/dashboard_venv/bin/python3 -u /opt/cryptomaster/start_flask_dashboard.py
Restart=always
RestartSec=10
# (start-rate limiting is disabled via StartLimitIntervalSec=0 in [Unit] above,
# so the API always keeps retrying and never stays down after a crash burst.)

# ── Systemd hardening (audit 9.6). Dashboard only READS data — no writes, no
# shell, no root, no SSH/Firebase key access. Each directive is compatible with
# a read-only Python+SQLite (mode=ro) service; verify with a service test before
# tightening further.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077
# The service opens cache.sqlite read-only; no path under /opt needs write.
ReadOnlyPaths=/opt/cryptomaster /opt/dashboard_venv

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dashboard-web

[Install]
WantedBy=multi-user.target
